# RAID1AutoRecovery/etc/init.d/lsrrbd
### BEGIN INIT INFO
# Provides:          md_auto_resync
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

#!/bin/sh
# (c) Copyright [2017] Hewlett Packard Enterprise Development LP
#
# This program is free software; you can redistribute it and/or modify it
# under the terms of the GNU General Public License, version 2, as
# published by the Free Software Foundation; either version 2 of the
# License.
#
# This program is distributed in the hope that it will be useful, but
# WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the GNU General
# Public License for more details.
#
# You should have received a copy of the GNU General Public License along
# with this program; if not, write to the Free Software Foundation, Inc.,
# 59 Temple Place, Suite 330, Boston, MA 02111-1307 USA
#
# Rev 1.1 04/26/2017

N=/etc/init.d/lsrrbd

case "$1" in
  start)
    /opt/hpe/lsrrb/bin/md_auto_resync.py &
    /opt/hpe/lsrrb/bin/HPEtemp.sh &
    ;;
  stop|reload|restart|force-reload|status)
    ;;
  *)
    echo "Usage: $N {stop|reload|restart|force-reload|status}"
    exit 1
    ;;
esac

exit 0
